# systemd/LG_Buddy_update_check.service
[Unit]
Description=LG Buddy background update check
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=@LG_BUDDY_BIN@ updates background-check
