# systemd/LG_Buddy_update_check.timer
[Unit]
Description=LG Buddy background update check timer
After=graphical-session.target
PartOf=graphical-session.target

[Timer]
OnCalendar=weekly
RandomizedDelaySec=12h
Persistent=true

[Install]
WantedBy=graphical-session.target
