# systemd/LG_Buddy_update_check.timer
[Unit]
Description=LG Buddy background update check timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
